# box/usbguard/usbguard.service.d/10-ca-usb-early.conf
# Systemd drop-in override for the vendor usbguard.service unit ("order
# usbguard.service early").
#
# The Debian usbguard package's shipped unit only has `[Install]
# WantedBy=basic.target` -- which pulls usbguard.service in as part of
# reaching basic.target (an early boot target, well before multi-user.target/
# getty/network), but WantedBy= alone does NOT order it relative to that
# target's completion (systemd.special(7): a unit that wants to be UP by the
# time a target is considered reached needs BOTH WantedBy=<target> AND
# Before=<target> -- the vendor unit has only the first). Without the
# matching Before=, systemd is free to start usbguard.service concurrently
# with (or, in principle, arbitrarily later than) the rest of basic.target's
# other dependencies.
#
# Adding the missing Before=basic.target here (a drop-in, NOT an edit to the
# vendor-owned /usr/lib/systemd/system/usbguard.service, so a package
# upgrade never clobbers it) makes explicit what "order early" means in
# practice: minimize the pre-usbguard-service boot window during which a
# USB device could enumerate before USBGuard is watching for it (see
# box/usbguard/README-authorized-default.md's discussion of that same
# window and why closing it further via a kernel cmdline parameter is
# comparatively low-value once this ordering fix is in place).
#
# NOTE: this file taking effect at boot is orthogonal to whether
# usbguard.service is *enabled* -- box/install.sh installs this override but
# deliberately leaves usbguard.service disabled/stopped; enabling it is a
# separate, deliberate operator step.
[Unit]
Before=basic.target
